# conf/wifi-check.service
[Unit]
Description=Verify Mediatek Wi-Fi modules
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/wifi-check.sh

[Install]
WantedBy=multi-user.target
